# deploy/kriegspiel-gpt-nano-bot.service
[Unit]
Description=Kriegspiel GPT Nano Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fil
WorkingDirectory=/home/fil/dev/kriegspiel/bot-gpt-nano
ExecStart=/home/fil/dev/kriegspiel/bot-gpt-nano/.venv/bin/python /home/fil/dev/kriegspiel/bot-gpt-nano/bot.py --poll-seconds 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
